# etc/profile.d/netrunner.sh
XDG_DATA_DIRS=/usr/share/netrunner-default-settings/usr-share:/usr/share:/usr/share/plasma:/usr/local/share 
